# config/supervisord.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm7 -F
autorestart=true
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
startretries=0